# config/keymon.service
[Unit]
Description=Key Monitor

[Service]
ExecStart=/usr/local/bin/keymon
Type=notify
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
