# config/keymon.service
[Unit]
Description=Key Monitor

[Service]
ExecStart=/usr/local/bin/keymon /dev/input/by-id/usb-Microsoft_Surface_Type_Cover-event-kbd
Type=notify
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
